# PaginaWeb/practica2.conf
[program:urban-train]
command=/home/ubuntu/web2/bin/gunicorn_start
user=jenkins
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/home/ubuntu/web2/logs/gunicorn-error.log